# obsidian_livesync/rootfs/etc/services.d/couchdb/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

export COUCHDB_USER="$(bashio::config 'COUCHDB_USER')"
export COUCHDB_PASSWORD="$(bashio::config 'COUCHDB_PASSWORD')"

if bashio::config.has_value 'COUCHDB_SECRET'; then
    export COUCHDB_SECRET="$(bashio::config 'COUCHDB_SECRET')"
fi

bashio::log.info "Starting CouchDB for Obsidian LiveSync..."
exec /docker-entrypoint.sh /opt/couchdb/bin/couchdb
